# prds/loop-revival-v2/deploy/openclaw/systemd/coding-agent-runner.timer
[Unit]
Description=Run coding-agent daily (last pipeline stage; only acts on status:ready-for-code)

[Timer]
# Daily at 05:50 — after engineering 05:10 + qa 05:30. Daily NOT hourly (the #42 runaway lesson).
OnCalendar=*-*-* *:13/15:00
Persistent=true
RandomizedDelaySec=90
Unit=coding-agent-runner.service

[Install]
WantedBy=timers.target
