# prds/loop-revival-v2/deploy/openclaw/systemd/coding-agent-runner.timer
[Unit]
Description=Run coding-agent ~every 15 min (backstop; webhook fires on status:ready-for-code)

[Timer]
OnCalendar=*-*-* 02/3:05:00
Persistent=true
RandomizedDelaySec=90
Unit=coding-agent-runner.service

[Install]
WantedBy=timers.target
