# prds/loop-revival-v2/deploy/openclaw/systemd/qa-runner.timer
[Unit]
Description=Run qa per-feature test-plan author daily (after engineering 05:10, before coding-agent 05:50)

[Timer]
# Daily at 05:30 — second stage of the pipeline. Daily NOT hourly (the #42 runaway lesson).
OnCalendar=*-*-* *:08/15:00
Persistent=true
RandomizedDelaySec=90
Unit=qa-runner.service

[Install]
WantedBy=timers.target
